# supercapfreezer.service
[Unit]
Description=SUPERCAPFREEZER Temperature Monitor
After=network.target
StartLimitInterval=60
StartLimitBurst=3

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/supercapfreezer
ExecStart=/home/pi/supercapfreezer/venv/bin/python main.py --port /dev/ttyACM0
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security: run as limited user, read-only root filesystem (if possible)
PrivateTmp=yes
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
