# supercapfreezer.service
[Unit]
Description=SUPERCAPFREEZER Temperature Monitor
After=network.target
StartLimitInterval=60
StartLimitBurst=3

[Service]
Type=simple
# Replace with your runtime user and install path before enabling.
User=REPLACE_WITH_USERNAME
WorkingDirectory=/home/REPLACE_WITH_USERNAME/supercapfreezer
ExecStart=/home/REPLACE_WITH_USERNAME/supercapfreezer/venv/bin/python /home/REPLACE_WITH_USERNAME/supercapfreezer/main.py --config /home/REPLACE_WITH_USERNAME/supercapfreezer/config.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

# Security: run as limited user, read-only root filesystem (if possible)
PrivateTmp=yes
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
